# publicbox/rpi/udev/99-wifi.rules
KERNEL=="wlan*", RUN+="/bin/sh -c /opt/publicbox/rpi/bin/starter.sh"
